# ubuntu_safe-shutdown.service
[Unit]
Description=Nespi4 Case Safe Shutdown

[Service]
Type=idle
ExecStart=PYTHON_EXEC /opt/RetroFlag/SafeShutdown.py
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=safe-shutdown.service
